# .config/systemd/user/foot.service
[Service]
ExecStart=/usr/bin/foot --server
Restart=always

[Install]
WantedBy=multi-user.target
