# .config/systemd/user/sysupdate-compile.service
[Service]
Type=oneshot
ExecStart=/home/radu/bin/sysupdate_compile.sh
