# scripts/.config/systemd/user/sysupdate-compile.timer
[Unit]
Description=Compile & update AUR/Emacs when idle twice a day (4 AM/PM)

[Timer]
OnCalendar=4/12:0

[Install]
WantedBy=timers.target
