# notmuch/.config/systemd/user/notmuch-new.service
[Service]
Type=oneshot
ExecStart=/usr/bin/notmuch new
TimeoutSec=600
Environment=IN_NOTMUCH_SYSTEMD=1
